# emmy_serve/systemd/emmy-sidecar.service
[Unit]
Description=Emmy sidecar — vLLM lifecycle controller (Phase 04.2)
Documentation=file:///%h/code/emmy/docs/runbook.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/code/emmy
ExecStart=%h/.local/bin/uv run python -m emmy_serve.swap.controller
Restart=on-failure
RestartSec=5
TimeoutStopSec=40
KillSignal=SIGTERM
Environment=PYTHONUNBUFFERED=1
Environment=VLLM_NO_USAGE_STATS=1
Environment=VLLM_LOAD_FORMAT=fastsafetensors
Environment=HF_HOME=%h/.cache/huggingface
# Phase 04.2 follow-up — the sidecar spawns subprocesses (orchestrator, docker)
# that need `uv` on PATH. systemd's user units inherit a minimal PATH that does
# NOT include ~/.local/bin (where uv lives by default), so without this the
# orchestrator launch dies with [Errno 2] No such file or directory and /start
# silently fails (the SSE stream emits an error frame but the slash command's
# default exit=0 gate masks it as success). Include /usr/local/bin and
# /usr/sbin too — docker, nvidia-smi, and other subprocess targets live there.
Environment=PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
